# systemd/aileron-daemon.service
[Unit]
Description=Aileron local AI daemon
Documentation=https://github.com/aileron-project/aileron
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/aileron-daemon
Restart=on-failure
RestartSec=3

# The socket lives in $XDG_RUNTIME_DIR which systemd already exposes to user services.
Environment="XDG_RUNTIME_DIR=%t"

[Install]
WantedBy=graphical-session.target
